# overlay/etc/s6-overlay/s6-rc.d/smb/run
#!/bin/withenv sh

set -euo pipefail


#======================================================================================================================
# Run preflight checks
#======================================================================================================================

nu -c "use bf-samba run ; run preflight"


#======================================================================================================================
# Start Samba in foreground mode.
#======================================================================================================================

smbd --configfile=${BF_SAMBA_CONF} --debug-stdout --foreground
